# playbooks/roles/app/files/supervisor.conf
[program:gunicorn]
command=gunicorn hello:application -c /etc/hello/gunicorn.conf.py
directory=/var/www/hello/
autostart=true
autorestart=true
redirect_stderr=True
